# config/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
pidfile=/tmp/supervisord.pid
nodaemon=false

[program:deefuzzer]
command=deefuzzer /home/deefuzzer/mystation/deefuzzer.xml
